# config/status-control.service
[Unit]
Description=Status Screen - Control Server
After=network-online.target

[Service]
User=__STATUS_SCREEN_USER__
Environment=STATUS_SCREEN_DIR=__STATUS_SCREEN_DIR__
WorkingDirectory=__STATUS_SCREEN_DIR__
ExecStart=__STATUS_SCREEN_DIR__/.venv/bin/python __STATUS_SCREEN_DIR__/pi/control_server.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
