# config/status-from-ics.service
[Unit]
Description=Status Screen - ICS Calendar Resolver
After=network-online.target

[Service]
User=__STATUS_SCREEN_USER__
Environment=STATUS_SCREEN_DIR=__STATUS_SCREEN_DIR__
WorkingDirectory=__STATUS_SCREEN_DIR__
ExecStart=__STATUS_SCREEN_DIR__/.venv/bin/python __STATUS_SCREEN_DIR__/pi/status_from_ics.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
